# packages/core/services/vfr-gather.service
[Unit]
Description=Video Face Recognition Gather Service

[Service]
Type=simple
User=mil
Group=mil
WorkingDirectory=/home/mil/vfr/packages/core
EnvironmentFile=/home/mil/vfr/packages/core/services/vfr-gather.conf
ExecStart=/usr/bin/xinit -e /bin/bash -c "yarn gather /home/mil/vfr/.data/training-images $NAME" $* -- :2

[Install]
WantedBy=multi-user.target
