# packages/core/services/vfr-train.service
[Unit]
Description=Video Face Recognition Train Service

[Service]
Type=simple
User=mil
Group=mil
WorkingDirectory=/home/mil/vfr/packages/core
EnvironmentFile=/home/mil/vfr/packages/core/services/vfr.conf
ExecStart=/bin/bash -c "yarn train /home/mil/vfr/.data/training-images /home/mil/vfr/.data/model.json"

[Install]
WantedBy=multi-user.target