# packages/core/services/vfr-recognize.service
[Unit]
Description=Video Face Recognition Recognize Service

[Service]
Type=simple
User=mil
Group=mil
WorkingDirectory=/home/mil/vfr/packages/core
EnvironmentFile=/home/mil/vfr/packages/core/services/vfr.conf
ExecStart=/bin/bash -c "yarn recognize /home/mil/vfr/.data/model.json"
Restart=always

[Install]
WantedBy=multi-user.target